# etc/udev/rules.d/99-power.rules
ACTION=="change", SUBSYSTEM=="power_supply", ATTRS{type}=="Mains", ATTRS{online}=="0", RUN+="/home/sharwesh/.local/bin/power_bat"
ACTION=="change", SUBSYSTEM=="power_supply", ATTRS{type}=="Mains", ATTRS{online}=="1", RUN+="/home/sharwesh/.local/bin/power_ac"
